# etc/profile.d/aliases.sh
alias psx='ps axfo pid,euser,egroup,args'
